# agent/systemd/agent-client.service
[Unit]
Description=Agentic Sandbox Agent Client
Documentation=https://git.integrolabs.net/roctinam/agentic-sandbox
After=network-online.target qemu-guest-agent.service
Wants=network-online.target

[Service]
Type=simple
User=agent
Group=agent
WorkingDirectory=/opt/agentic-sandbox/agent

# Environment file with agent credentials
EnvironmentFile=/etc/agentic-sandbox/agent.env

# Python agent client
ExecStart=/usr/bin/python3 /opt/agentic-sandbox/agent/grpc_client.py

# Restart configuration
Restart=always
RestartSec=5
StartLimitBurst=5
StartLimitIntervalSec=300

# Resource limits
MemoryMax=512M
CPUQuota=50%

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
ReadWritePaths=/mnt/inbox
ReadOnlyPaths=/mnt/global /etc/agentic-sandbox

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=agent-client

[Install]
WantedBy=multi-user.target
