# agent-rs/systemd/agent-client.service
[Unit]
Description=Agentic Sandbox Agent Client (Rust)
Documentation=https://git.integrolabs.net/roctinam/agentic-sandbox
After=network-online.target qemu-guest-agent.service
Wants=network-online.target

[Service]
Type=simple
User=agent
Group=agent
WorkingDirectory=/opt/agentic-sandbox

# Environment file with agent credentials
EnvironmentFile=/etc/agentic-sandbox/agent.env

# Rust agent client (compiled binary)
ExecStart=/opt/agentic-sandbox/bin/agent-client

# Restart configuration
Restart=always
RestartSec=5
StartLimitBurst=5
StartLimitIntervalSec=300

# VM-level isolation is the security boundary.
# The agent (and its PTY sessions) needs full OS access to install
# software, run sudo, and manage the system on behalf of the user.
# Resource limits are enforced by KVM (vcpu/memory), not systemd.

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=agent-client

[Install]
WantedBy=multi-user.target
