# deploy/packaging/systemd/agent-client.service
[Unit]
Description=Agentic Sandbox Agent Client
Documentation=https://github.com/jmagly/agentic-sandbox
After=network-online.target
Wants=network-online.target
StartLimitBurst=5
StartLimitIntervalSec=300

[Service]
Type=simple
EnvironmentFile=-/etc/agentic-sandbox/agent.env
ExecStart=/usr/bin/agent-client
Restart=always
RestartSec=5
TimeoutStopSec=30
KillMode=process
LimitNOFILE=65536
StandardOutput=journal
StandardError=journal
SyslogIdentifier=agent-client

[Install]
WantedBy=multi-user.target
